# files/etc/cron.d/dsa-puppet-stuff
@hourly  root [ ! -d /var/cache/dsa ] || touch /var/cache/dsa/cron.alive
