# files/etc/cron.d/dsa-puppet-stuff
SHELL=/bin/bash
@hourly  root [ ! -d /var/cache/dsa ] || touch /var/cache/dsa/cron.alive
34 */4 * * * root if [ -x /usr/sbin/puppetd ]; then sleep $(( $RANDOM \% 3600 )); if [ -x /usr/bin/timeout ]; then TO="timeout 3600"; else TO=""; fi; $TO /usr/sbin/puppetd --factsync -o --no-daemonize 2>&1 | grep -v 'v6: error fetching interface information: Device not found' ; fi
